# cpack/productbuild/ChmodBPF/ChmodBPF.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>@CU_COMPANY_DOMAIN@.@CU_PROJECT_COMPANYNAME@.@CMAKE_INSTALL_DEFAULT_COMPONENT_NAME@.ChmodBPF</string>
	<key>RunAtLoad</key>
	<true/>
	<key>Program</key>
	<string>/Library/Application Support/@CU_PROJECT_COMPANYNAME@/@PROJECT_NAME@/ChmodBPF/ChmodBPF</string>
</dict>
</plist>
